# scripts/launchd/com.nexus.web.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.nexus.web</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/shadow/code/nexus/scripts/launchd/with-node.sh</string>
    <string>/Users/shadow/code/nexus/nexus-web</string>
    <string>npm</string>
    <string>run</string>
    <string>dev</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/shadow/code/nexus/nexus-web</string>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <dict>
    <key>SuccessfulExit</key>
    <false/>
    <key>Crashed</key>
    <true/>
  </dict>

  <key>ThrottleInterval</key>
  <integer>10</integer>

  <key>StandardOutPath</key>
  <string>/Users/shadow/Library/Logs/Nexus/web.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/shadow/Library/Logs/Nexus/web.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>NODE_ENV</key>
    <string>development</string>
  </dict>

  <key>ProcessType</key>
  <string>Interactive</string>
</dict>
</plist>
