# scripts/launchd/com.nexus.arena.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.nexus.arena</string>

  <key>ProgramArguments</key>
  <array>
    <string>__REPO_ROOT__/scripts/launchd/with-node.sh</string>
    <string>__REPO_ROOT__/arena</string>
    <string>npm</string>
    <string>run</string>
    <string>dev</string>
  </array>

  <key>WorkingDirectory</key>
  <string>__REPO_ROOT__/arena</string>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <dict>
    <key>SuccessfulExit</key>
    <false/>
    <key>Crashed</key>
    <true/>
  </dict>

  <key>ThrottleInterval</key>
  <integer>10</integer>

  <key>StandardOutPath</key>
  <string>__LOG_DIR__/arena.log</string>

  <key>StandardErrorPath</key>
  <string>__LOG_DIR__/arena.log</string>

  <key>ProcessType</key>
  <string>Interactive</string>
</dict>
</plist>
